# supervisord.conf
[supervisord]
nodaemon=true

[program:websockify]
command=websockify --web /usr/share/novnc 8080 %(ENV_VNC_SERVER)s
autorestart=true

